# salt/afval/systemd.timer
[Unit]
Description=Afvalkalendar Timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=8h

[Install]
WantedBy=timers.target
